# providers/systemd/_test/systemd-test-fail.service
[Unit]
Description=Systemd provider test - Fail

[Service]
ExecStart=/bin/sh -c "exit 1"

[Install]
WantedBy=multi-user.target