# providers/systemd/_test_data/systemd-test-exit1.service
[Unit]
Description=Systemd provider test - Exit 1

[Service]
ExecStart=/bin/sh -c "exit 1"

[Install]
WantedBy=multi-user.target